# ansible-inventory/playbooks/files/prom-rpi-exporter.service
[Unit]
Description=Collect Raspberry Pi metrics for Prometheus

[Service]
Type=oneshot
ExecStart=/usr/local/bin/prom-rpi-exporter.sh
