# ansible-inventory/playbooks/files/prom-rpi-exporter.timer
[Unit]
Description=Raspberry Pi prometheus metrics timer

[Timer]
OnCalendar=*:*:0

[Install]
WantedBy=timers.target
